# roles/roshi-server/files/roshi-upstart.conf
description "roshi server"

start on started redis-server
stop on shutdown

exec sudo -u roshi  /usr/local/sbin/roshi-server -redis.instances=localhost:6379

respawn
respawn limit 10 5
